# modules/firstboot/files/usr/lib/systemd/user/firstboot-login.service
[Unit]
Description=Run custom first login scripts
After=graphical-session.target
ConditionPathExists=!%h/.firstboot-login-done

[Service]
Type=oneshot
ExecStart=/usr/lib/firstboot-scripts/run-login.sh
RemainAfterExit=no

[Install]
WantedBy=default.target
